# etc/systemd/system/cache-mover.service
[Unit]
Description=Move old files from SSD to HDD storage
After=mnt-storage.mount
Conflicts=snapraid-sync.service

[Service]
Type=oneshot
# Move files older than 30 days
ExecStart=/bin/bash -c 'cd /mnt/cache && find . -type f -mtime +30 -print | rsync -av --remove-source-files --files-from=- . /mnt/storage/'
# Clean up empty directories
ExecStart=/usr/bin/find /mnt/cache -type d -empty -delete
User=root
Nice=19
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
