# etc/systemd/system/snapraid-sync.timer
[Unit]
Description=Run SnapRAID sync daily

[Timer]
OnCalendar=daily
AccuracySec=1hour
Persistent=true

[Install]
WantedBy=timers.target
